# systemd/user/random-wallpaper.service
[Service]
ExecStart=%h/.config/scripts/set_random %h/Pictures/Backgrounds/
